# config/datasaver-lid.service
[Unit]
Description=Monitor laptop lid and toggle datasaver
After=multi-user.target

[Service]
Type=simple
ExecStart=/usr/bin/bash /opt/datasaver/lid_monitor.sh
Restart=always

[Install]
WantedBy=multi-user.target
