# datasaver/datasaver.service
[Unit]
Description=Datasaver GUI Application with Lid Watcher
After=network.target graphical.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 /usr/local/bin/datasaver-gui.py
Restart=on-failure
User=root
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/YOUR_USERNAME/.Xauthority

[Install]
WantedBy=default.target
